# XDG/.config/autostart/ibus.desktop
[Desktop Entry]
Name=LXTerminal
Type=Application
Exec=ibus-daemon -dr
Terminal=false

